# deploy/morfocampo-update.timer
[Unit]
Description=Periodic Morfocampo guarded update check

[Timer]
OnBootSec=10min
OnUnitActiveSec=1d
RandomizedDelaySec=30min
Persistent=true

[Install]
WantedBy=timers.target
